# CS183 Lab2/crontab
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root

*/2 * * * * root /root/cs183lab/lab2_part3.sh